# src/agents/flyer/systemd/flyer-recovery-watchdog.timer
[Unit]
Description=Flyer Studio recovery watchdog timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=flyer-recovery-watchdog.service

[Install]
WantedBy=timers.target
