# src/agents/expense_bookkeeper/systemd/prune-expense-receipts.service
[Unit]
Description=Agent #21 Expense Bookkeeper - prune expired receipts + expire stale leads
After=network.target

[Service]
Type=oneshot
# Agent #21 is optional on some VPSes. ExecCondition skips cleanly when the
# customer config is absent, instead of failing the timer.
ExecCondition=/usr/bin/test -f /opt/shift-agent/config.yaml
ExecCondition=/usr/bin/test -x /usr/local/lib/hermes-agent/venv/bin/python
ExecStart=/usr/local/lib/hermes-agent/venv/bin/python /usr/local/bin/prune-and-expire-expenses.py
User=shift-agent
Group=shift-agent
StandardOutput=append:/opt/shift-agent/logs/prune-expense.log
StandardError=append:/opt/shift-agent/logs/prune-expense.log
