# src/agents/expense_bookkeeper/systemd/prune-expense-receipts.service
[Unit]
Description=Agent #21 Expense Bookkeeper — prune expired receipts + expire stale leads
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/shift-agent/venv/bin/python /usr/local/bin/prune-and-expire-expenses.py
User=shiftagent
Group=shiftagent
StandardOutput=append:/var/log/shift-agent/prune-expense.log
StandardError=append:/var/log/shift-agent/prune-expense.log
