# src/agents/shift/systemd/shift-agent-health-watchdog.service
[Unit]
Description=Shift Agent — watchdog for the health check itself

[Service]
Type=oneshot
User=shift-agent
Group=shift-agent
EnvironmentFile=/opt/shift-agent/.env
ExecStart=/usr/local/bin/shift-agent-health-watchdog.sh
