# src/agents/catering/systemd/catering-dispatcher-watchdog.service
[Unit]
Description=Catering Dispatcher Watchdog — F7 absorbing shim for missed catering SKILL invocation (Hermes WhatsApp lacks auto_skill binding)
After=network-online.target hermes-gateway.service
Wants=network-online.target
PartOf=hermes-gateway.service

[Service]
Type=simple
User=shift-agent
Group=shift-agent
EnvironmentFile=/opt/shift-agent/.env
Environment=HOME=/opt/shift-agent
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Long-running tail-style daemon. Restart on crash with backoff.
ExecStart=/root/.hermes/hermes-agent/venv/bin/python /usr/local/bin/catering-dispatcher-watchdog
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=10

# Stdout/stderr to file (matches gateway pattern)
StandardOutput=append:/opt/shift-agent/logs/catering-dispatcher-watchdog.log
StandardError=append:/opt/shift-agent/logs/catering-dispatcher-watchdog.log

# Security hardening (mirrors send-routing-accuracy-summary pattern)
NoNewPrivileges=true
ProtectSystem=strict
# Daemon needs to read bridge.log (in /root/.hermes/, not writable to shift-agent),
# write to /opt/shift-agent/logs/, and invoke /usr/local/bin/* scripts.
ReadOnlyPaths=/root/.hermes
ReadWritePaths=/opt/shift-agent
ProtectHome=read-only
PrivateTmp=true
RuntimeDirectory=shift-agent
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
