# src/platform/systemd/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway - WhatsApp + skill runtime for Shift Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=shift-agent
Group=shift-agent
Environment=HOME=/opt/shift-agent
Environment=HERMES_HOME=/root/.hermes
# PATH must include /root/.hermes/node/bin so Hermes can spawn the WhatsApp
# bridge subprocess (`["node", "bridge.js"]`). systemd's default service PATH
# lacks this dir; without the override, Hermes logs "Node.js not installed
# or bridge not configured" and refuses to start the WhatsApp adapter.
Environment=PATH=/root/.hermes/node/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=FLYER_ALLOW_INTEGRATED_POSTER=1
EnvironmentFile=/opt/shift-agent/.env

# Targeted permissions guard. Runs with elevated privileges even though the
# gateway itself runs as shift-agent. Avoid broad recursive chown of all
# /root/.hermes because stale operator backup files can block gateway startup.
ExecStartPre=+/usr/local/bin/shift-agent-hermes-permissions

ExecStart=/root/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main gateway run --replace
Restart=on-failure
RestartSec=30
StandardOutput=append:/opt/shift-agent/logs/hermes-gateway.log
StandardError=append:/opt/shift-agent/logs/hermes-gateway.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/shift-agent /root/.hermes
# ProtectHome=true would hide /root even with ReadWritePaths exception — avoid.
ProtectHome=read-only
PrivateTmp=true
RuntimeDirectory=shift-agent

[Install]
WantedBy=multi-user.target
