# src/platform/systemd/triveni-portal.service
[Unit]
Description=Triveni SMB-Agents Portfolio Portal (static HTML)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/triveni/portal/index.html

[Service]
Type=simple
User=shift-agent
Group=shift-agent
ExecStart=/usr/bin/python3 -m http.server 8080 --directory /opt/triveni/portal --bind 0.0.0.0
Restart=on-failure
RestartSec=30
StartLimitBurst=3
StandardOutput=append:/opt/shift-agent/logs/triveni-portal.log
StandardError=append:/opt/shift-agent/logs/triveni-portal.log

# Security hardening (mirrors hermes-gateway.service convention).
# R2-B1 fix from design review: ReadWritePaths required for StandardOutput
# append target. Without it, ProtectSystem=strict mounts /opt read-only and
# the service fails at startup with EROFS on triveni-portal.log.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/shift-agent/logs
# R2-M1 fix: read-only matches deployed convention. ProtectHome=true would
# hide /root even with explicit ReadWritePaths exception (see comment in
# hermes-gateway.service:33-34).
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
