# src/agents/compliance/systemd/check-compliance-deadlines.service
[Unit]
Description=Compliance Calendar daily reminder check (Agent #13)
After=network-online.target hermes-gateway.service
Wants=network-online.target
ConditionFileIsExecutable=/usr/local/bin/check-compliance-deadlines.py
ConditionPathExists=/opt/shift-agent/config.yaml

[Service]
Type=oneshot
User=shift-agent
Group=shift-agent
EnvironmentFile=/opt/shift-agent/.env
Environment=HOME=/opt/shift-agent
# PR-Agent13-v0.1 (Reviewer A2-B2): explicit TZ for OnCalendar — Hetzner VPS
# default tz is UTC; without this, OnCalendar=*-*-* 06:00:00 fires at 06:00 UTC
# (= 01:00 EST / 22:00-prior-day PST). The deploy script substitutes
# @CUSTOMER_TZ@ from cfg.customer.timezone at install time via:
#   python3 -c "import yaml; print(yaml.safe_load(open('/opt/shift-agent/config.yaml'))['customer']['timezone'])"
Environment=TZ=@CUSTOMER_TZ@
ExecStart=/usr/local/bin/check-compliance-deadlines.py
StandardOutput=journal
StandardError=journal

# Security hardening (mirrors send-daily-brief.service:18-23)
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/shift-agent
ProtectHome=read-only
PrivateTmp=true
RuntimeDirectory=shift-agent

[Install]
WantedBy=multi-user.target
