# src/agents/daily_brief/systemd/send-daily-brief.timer
[Unit]
Description=Daily Brief timer (15-min poll, self-gates on cfg.daily_brief.brief_time)
Requires=send-daily-brief.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
AccuracySec=1min
Unit=send-daily-brief.service

[Install]
WantedBy=timers.target
